# jobs/crontab
# [second(s)] minute(s) hour(s) day(s)_of_month month(s) day(s)_of_week command

11 * * * * /root/update_top_trending.sh
5 * * * * /root/MetaMaint.sh
