# jobs/crontab
# [second(s)] minute(s) hour(s) day(s)_of_month month(s) day(s)_of_week command

11 * * * * /root/update_top_trending.sh
13 3 * * 1 /root/MetaMaint.sh
47 3 * * * /root/update_spam_blacklist.sh
